# docker/web/crontab
# Executa a todo minuto
* * * * * php /var/www/html/teste-cron.php
